# file-bkp/debian/rsync-conf/etc/cron.d/rsync-conf
# /etc/cron.d/rsync: crontab entries for the rsync package
#
# Backup will run twice a day
#
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/etc/rsync

57 0/12 * * * ubuntu /etc/rsync/backup.sh
