# file-bkp/inc_rst.list
#please list here dirs and files for file-restore
#for example:
#gw server:
/home/ubuntu/.ssh
/etc/file-bkp
/etc/cron.d/file-bkp
/etc/iptables/rules.v4
/etc/openvpn
