# prometheus-node-exporter/debian/prometheus-node-exporter/DEBIAN/postinst
#!/bin/sh
# postinst script for prometheus-node-exporter
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)


    if ! dpkg-statoverride --list /var/lib/prometheus/node-exporter >/dev/null ; then
        chown prometheus:prometheus /var/lib/prometheus/node-exporter
    fi

#    if ! dpkg-statoverride --list /etc/prometheus >/dev/null ; then
#        chown prometheus:prometheus /etc/prometheus
#    fi
#    chmod 2750 /var/log/prometheus
#    chmod 2755 /usr/bin/prometheus
#    chmod 2755 /usr/bin/promtool
#    chmod 2775 /var/lib/prometheus
#    chown -R :prometheus /usr/bin/prometheus
#    chown -R :prometheus /usr/bin/promtool
#    chown -R :prometheus /var/log/prometheus
#    chown -R :prometheus /var/lib/prometheus

#    systemctl restart prometheus
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installsystemd/12.10ubuntu1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# This will only remove masks created by d-s-h on package removal.
	deb-systemd-helper unmask 'prometheus-node-exporter.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'prometheus-node-exporter.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'prometheus-node-exporter.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'prometheus-node-exporter.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/12.10ubuntu1
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'prometheus-node-exporter.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
